# dist/linux/newrelic/init/nr-counterpartyd.conf
#!upstart
description "counterpartyd (mainnet) (newrelic)"
limit nofile 16384 16384

# When to start the service
start on startup

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn
respawn limit 10 60 # 10 times in 60 seconds

script
USER=xcpd
USER_HOME=/home/xcp
NR_ADMIN_PATH=$USER_HOME/counterpartyd_build/env/bin/newrelic-admin
PYTHON_INTERP_PATH=$USER_HOME/counterpartyd_build/env/bin/python3
PYTHON_SCRIPT_PATH=$USER_HOME/counterpartyd_build/dist/counterpartyd/counterpartyd.py
NEW_RELIC_CONFIG_FILE=/etc/newrelic/nr_counterpartyd.ini
export NEW_RELIC_CONFIG_FILE

start-stop-daemon --start --pidfile $USER_HOME/.config/counterpartyd/counterpartyd.pid --chuid $USER --exec $NR_ADMIN_PATH -- run-program $PYTHON_INTERP_PATH $PYTHON_SCRIPT_PATH server
end script
