# dist/linux/newrelic/init/nr-counterblockd.conf
#!upstart
description "counterblockd (mainnet) (newrelic)"

# When to start the service
start on startup

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn
respawn limit 10 60 # 10 times in 60 seconds

script
USER=xcp
USER_HOME=$(eval echo ~${USER})
NR_ADMIN_PATH=$USER_HOME/counterpartyd_build/env.cwalletd/bin/newrelic-admin
PYTHON_INTERP_PATH=$USER_HOME/counterpartyd_build/env.cwalletd/bin/python
PYTHON_SCRIPT_PATH=$USER_HOME/counterpartyd_build/dist/counterblockd/counterblockd.py
NEW_RELIC_CONFIG_FILE=/etc/newrelic/nr_counterblockd.ini
export NEW_RELIC_CONFIG_FILE

start-stop-daemon --start --pidfile $USER_HOME/.config/counterblockd/counterblockd.pid --chuid $USER --exec $NR_ADMIN_PATH -- run-program $PYTHON_INTERP_PATH $PYTHON_SCRIPT_PATH
end script
