# dist/linux/runit/counterblockd-testnet/run
#!/bin/sh

USER=xcpd
USER_HOME=/home/xcp

exec chpst -u ${USER} /usr/local/bin/counterblockd --data-dir=${USER_HOME}/.config/counterblockd-testnet --testnet
