# dist/linux/runit/armory_utxsvr-testnet/run
#!/bin/sh

USER=xcpd
USER_HOME=/home/xcp

#sv start bitcoind-testnet || exit 1
exec 2>&1
exec chpst -u ${USER} /usr/local/bin/armory_utxsvr --testnet
