# dist/linux/runit/counterblock-testnet/run
#!/bin/sh

USER=xcpd
export USER_HOME=/home/xcp
ulimit -n 32768
sleep 2 #prevent cpu maxing out if cbd has a crashing problem
exec su -s /bin/bash -c 'PYTHON_EGG_CACHE=/tmp/counterblock.python-eggs /usr/local/bin/counterblock --testnet --config-file=${USER_HOME}/.config/counterblock/server.testnet.conf' ${USER}
