# dist/linux/init/bitcoind-testnet.conf.template
#!upstart
description "bitcoind (testnet)"

start on startup
stop on shutdown
oom never

respawn
respawn limit 10 60 # 10 times in 60 seconds

post-stop script
   /bin/sleep 5
end script

script
USER=!RUN_AS_USER!
USER_HOME=!USER_HOMEDIR!
start-stop-daemon --start -c $USER --make-pidfile --pidfile $USER_HOME/.bitcoin-testnet/bitcoind.pid --exec /usr/bin/bitcoind -- -daemon=0 -testnet -datadir=$USER_HOME/.bitcoin-testnet
end script