# dist/linux/init/bitcoind-testnet.conf.template
#!upstart
description "bitcoind (testnet)"

start on startup
stop on shutdown
oom never
respawn
respawn limit 10 60 # 10 times in 60 seconds

env USER=!RUN_AS_USER!
exec start-stop-daemon --start -c $USER --pidfile ~$USER/.bitcoin-testnet/bitcoind.pid --startas /usr/bin/bitcoind -- -daemon=0 -testnet -datadir=~$USER/.bitcoin-testnet
