# dist/linux/runit/nginx/run
#!/bin/sh

USER=xcpd
USER_HOME=/home/xcp

exec chpst -u ${USER} /usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"