# dist/linux/runit/bitcoind/run
#!/bin/sh

USER=xcpd
USER_HOME=/home/xcp

exec chpst -u ${USER} bitcoind -daemon=0 -datadir=${USER_HOME}/.bitcoin