# dist/linux/runit/insight-testnet/run
#!/bin/sh

USER=xcpd
USER_HOME=/home/xcp
chown -R $USER $USER_HOME/insight-api
export INSIGHT_NETWORK=testnet
export INSIGHT_DB=$USER_HOME/insight-api/db
export BITCOIND_DATADIR=$USER_HOME/.bitcoin-testnet/
export BITCOIND_USER=`cat $USER_HOME/.bitcoin-testnet/bitcoin.conf | sed -n 's/.*rpcuser=\([^ \n]*\).*/\1/p'`
export BITCOIND_PASS=`cat $USER_HOME/.bitcoin-testnet/bitcoin.conf | sed -n 's/.*rpcpassword=\([^ \n]*\).*/\1/p'`
#BITCOIND_HOST -- use the default (localhost)
#BITCOIND_PORT -- use the default (18332)
#BITCOIND_P2P_PORT -- use the default (18333)

#sv start bitcoind-testnet || exit 1
cd ${USER_HOME}/insight-api
exec 2>&1
exec chpst -u ${USER} /usr/bin/node ${USER_HOME}/insight-api/insight.js
