# dist/linux/runit/counterblock/run
#!/bin/sh

USER=xcpd
export USER_HOME=/home/xcp
ulimit -n 32768
sleep 2 #prevent cpu maxing out if cbd has a crashing problem
exec su -s /bin/bash -c '/usr/local/bin/counterblock --data-dir=${USER_HOME}/.config/counterblock' ${USER}
